# deploy/systemd/meepo-bot.service
[Unit]
Description=Meepo Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=meepo
Group=meepo
WorkingDirectory=/home/meepo/meepo-bot
EnvironmentFile=/etc/meepo/meepo-bot.env
ExecStart=/usr/bin/env npm run dev:bot
Restart=always
RestartSec=5
TimeoutStopSec=20
KillSignal=SIGTERM
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
